# uwsgi-cloud.conf
[program:uwsgi-cloud]
command=uwsgi --http :8000 --wsgi-file cloud/wsgi.py
autostart=true
directory=/home/ubuntu/cloud-series/cloud/
